# resources/fty-mdns-sd.service.in
[Unit]
Description=@PROJECT_NAME@ service: make this system discoverable via AVAHI
After=malamute.service network.target fty-info.service avahi-daemon.service
Requires=malamute.service network.target fty-info.service avahi-daemon.service
Conflicts=rescue.target shutdown.target poweroff.target halt.target reboot.target emergency.target

Requisite=bios-allowed.target
After=bios-allowed.target
PartOf=bios-allowed.target
BindsTo=bios-allowed.target
Conflicts=bios-shutdown.target

# fty-info/fty-mdns-sd should be pre-EULA, so that just-deployed instances
# are visible on the network for remote (distributed, through central point)
# commissioning, eventually.
PartOf=bios-pre-eula.target

[Service]
Type=simple
User=@AGENT_USER@
Restart=always

Environment='SYSTEMD_UNIT_FULLNAME=%n'
Environment="prefix=/usr"
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/usr/share/fty/etc/default/fty
EnvironmentFile=-/usr/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/fty
EnvironmentFile=-/etc/default/fty__%n.conf

ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/@PROJECT_NAME@ -c @AGENT_CONF_FILE@

[Install]
WantedBy=bios-pre-eula.target
