# resources/fty-mdns-sd.service.in
[Unit]
Description=fty-mdns-sd service: make this system discoverable via AVAHI
After=malamute.service network.target fty-info.service avahi-daemon.service
Requires=malamute.service network.target fty-info.service avahi-daemon.service
Conflicts=rescue.target shutdown.target poweroff.target halt.target reboot.target emergency.target

Requisite=bios-allowed.target
After=bios-allowed.target
PartOf=bios-allowed.target
BindsTo=bios-allowed.target
Conflicts=bios-shutdown.target

# fty-info/fty-mdns-sd should be pre-EULA, so that just-deployed instances
# are visible on the network for remote (distributed, through central point)
# commissioning, eventually.
PartOf=bios-pre-eula.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-mdns-sd -c @sysconfdir@/@PACKAGE@/fty-mdns-sd.cfg

[Install]
WantedBy=bios-pre-eula.target
